# etc/profile.d/00.basic.sh
#!/bin/bash

BLACK='\e[01;30m'
RED='\e[01;31m'
GREEN='\e[01;32m'
YELLOW='\e[01;33m'
BLUE='\e[01;34m'
MAGENTA='\e[01;35m'
CYAN='\e[01;36m'
WHITE='\e[01;37m'
NC='\e[00m'

export HISTSIZE=10000
export HISTFILESIZE=20000
export HISTCONTROL=ignoredups
export HISTTIMEFORMAT="%F %T "
