# etc/profile.d/01.aliases.sh
#!/bin/bash

alias l="ls -h"
alias lll="ls -alh"
alias ll="ls -alh"
alias c="clear"
alias vi="vim"
