# etc/profile.d/02.env.sh
#!/bin/bash

#export LESSCHARSET=utf-8
#export LANG=zh_CN.UTF-8

