# etc/profile.d/22.git2.sh
alias gco='git checkout'
alias gcb='git checkout -b'
alias gcp='git cherry-pick'
