# etc/profile.d/91.hello.sh
#!/bin/bash

#python3 ~/.shrc/bin/pycow/pycow.py -n

