# etc/profile.d/91.hello.sh
#!/bin/bash

python3 ~/.dotusr/bin/pycow/pycow.py -n

